# ansiblerepo/roles/layercake/templates/layercake-pod.service.j2
# layercake-pod.service
# autogenerated by Podman 4.1.0
# Mon Jun  6 22:20:52 UTC 2022

[Unit]
Description=Podman layercake-pod.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=
Requires=
Before=

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/layercake-pod.pid %t/layercake-pod.pod-id
ExecStartPre=/usr/bin/podman pod create \
    --infra-conmon-pidfile %t/layercake-pod.pid \
    --pod-id-file %t/layercake-pod.pod-id \
    --name lc{{ '%02d' | format(compute_node_number) }} \
    --share ipc,net,uts,pid \
    --network cvlan -p 6818:6818 \
    --ip {{ netconfig['cvlan']['pod_ipv4_address'] }} \
    --replace
ExecStart=/usr/bin/podman pod start --pod-id-file %t/layercake-pod.pod-id
ExecStop=/usr/bin/podman pod stop --ignore --pod-id-file %t/layercake-pod.pod-id -t 10
ExecStopPost=/usr/bin/podman pod rm --ignore -f --pod-id-file %t/layercake-pod.pod-id
PIDFile=%t/layercake-pod.pid
Type=forking

[Install]
WantedBy=default.target
