# ansiblerepo/roles/layercake/templates/layercake-munged.service.j2
# layercake-munged.service
# autogenerated by Podman 4.1.0
# Mon Jun  6 21:40:17 UTC 2022

[Unit]
Description=Podman layercake-munged.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=/bin/rm -f %t/%n.ctr-id
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	-d \
	--replace \
	-ti \
	--pod layercake \
	-v /tmp/layercakepodtmp:/tmp:z \
	-v /run/layercake/munge:/run/munge:z \
	--secret mungekey \
	--name munged head:5000/munged
ExecStop=/usr/bin/podman stop --ignore --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm -f --ignore --cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
